# kombucha-tracker/kombucha.service
[Unit]
Description=Kombucha Tracker Backend
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/kombucha/kombucha-tracker
ExecStart=/usr/bin/node server.js
Restart=always
Environment=PORT=3001
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
